# services/weather-capture.service
[Unit]
Description=Capture sensor values to MariaDB
After=multi-user.target
[Service]
Type=simple
User=pi
EnvironmentFile=/etc/environment
ExecStart=/usr/bin/python3 /home/pi/weather-station/capture.py
[Install]
WantedBy=multi-user.target